# dimecash/dc.crontab
#!/bin/bash


# * * * * * /root/dimecash/dc.crontab >> /root/dimecash/cron.log 2>&1

. /root/dimecash/dc.bash

cd /root/dimecash

LAST=`cat dc.last`

CURRENT=`dc.blockcount`

echo $LAST $CURRENT
